# community/deepin/live-overlay/etc/skel/.config/autostart/deepin-live.desktop
[Desktop Entry]
Comment=Deepin live session settings
Exec=deepin-live
Name=deepin-live
Type=Application
X-Deepin-Vendor=user-custom
