# community/webdad/desktop-overlay/etc/xdg/autostart/keyboard.desktop
[Desktop Entry]
Name=setxkbmap
Comment=setxkbmap
Exec=/bin/setKeyboard.sh
StartupNotify=false
Terminal=false
Type=Application
