# community/cutefish/desktop-overlay/etc/systemd/system/multi-user.target.wants/manjaro-tweaks.service
[Unit]
Description=Manjaro Tweaks
 
[Service]
ExecStart=/usr/bin/sh /etc/manjaro/manjaro-tweaks.sh
 
[Install]
WantedBy=multi-user.target
